# deploy/systemd/pgsim.service
[Unit]
Description=ISO 8583 Simulator (pgsim)
After=network.target

[Service]
Type=simple
User=pgsim
WorkingDirectory=/opt/pgsim
ExecStart=/usr/bin/java -jar /opt/pgsim/pgsim.jar
Restart=on-failure
RestartSec=5
Environment=SERVER_PORT=8081
Environment=SIMULATOR_TCP_PORT=8080
# Environment=PGSIM_TCP_TLS_ENABLED=true
# Environment=PGSIM_TCP_TLS_KEYSTORE=/opt/pgsim/tls/server.p12

[Install]
WantedBy=multi-user.target
